# systemd/nexus-api.service
[Unit]
Description=Nexus API — Distributed Agent Platform REST API
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=nexus
Group=nexus
WorkingDirectory=/opt/nexus/Services/nexus
ExecStart=/usr/bin/node packages/api/dist/index.js
EnvironmentFile=/opt/nexus/Services/nexus/.env
Environment=PATH=/opt/nexus/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexus-api

[Install]
WantedBy=multi-user.target
