# systemd/nexus-health-watchdog.timer
[Unit]
Description=Nexus Platform Health Watchdog — every 60s

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=nexus-health-watchdog.service

[Install]
WantedBy=timers.target
